# patch/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
echo 'usb2' | tee /sys/bus/usb/drivers/usb/unbind
exit 0
